# docker/pi-hole-wireguard/etc/systemd/system/docker-wireguard-pihole-sync.service
[Unit]
Description=Pihole Sync script service
Wants=docker-wireguard-pihole-sync.timer

[Service]
Type=oneshot
ExecStart=/opt/utils/pihole-sync-slave.sh

[Install]
WantedBy=multi-user.target
