# docker/pi-hole-wireguard/etc/cron.d/pihole-sync
#m	h	d	m	wd	usr	cmd
*/15	*	*	*	*	root	/opt/utils/pihole-sync-slave.sh
