# containerd/npm/etc/systemd/system/npm.service
[Unit]
Description= Nginx Proxy Manager - Container
Requires=containerd.service
After=containerd.service
[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/containerd/npm/
ExecStart=/usr/local/bin/nerdctl compose up -d
ExecStop=/usr/local/bin/nerdctl compose down
TimeoutStartSec=0
[Install]
WantedBy=multi-user.target
